# profile.d/always_in_a_screen_session.sh
#if [[ -n "${SSH_TTY}" ]] && [[ -z "${STY}" ]]; then
#   screen -xRR && exit
#fi