# profile.d/kcsdw.sh
#!/bin/bash
alias kcsdw="mysql -u kcsdw -h ibm-x3550m4-1.gsslab.rdu2.redhat.com -p kcsdw"