# profile.d/vim.sh
alias vi=/usr/bin/vimx
alias vim=/usr/bin/vimx