# autostart/devilspie.desktop
# This file was automatically generated for the autostart                       
# item devilspie
[Desktop Entry]
Type=Application
Exec=/home/dvf1976/bin/devilspie
Hidden=False
Terminal=False
StartupNotify=False
Version=0.9.4
Name=devilspie
